# docker/services/api/finish
#! /bin/bash
echo "api exited with code" $?
sleep 3
