# docker/services/api/run
#! /bin/bash
# exec node --inspect-brk=0.0.0.0:9229 /api/index.mjs
exec node /api/index.mjs
