# docker/services/api/run
#! /bin/bash
# ${Math.random()}
exec node /api/index.mjs
