# docker/services/pdns/finish
#! /bin/bash
echo "pdns exited with code" $?
sleep 3
