# docker/services/memhunt/finish
#! /bin/bash
echo "memhunt exited with code" $?
sleep 3
