# docker/services/rsync/finish
#! /bin/bash
echo "rsync exited with code" $?
sleep 3
