# docker/services/rsync/run
#! /bin/bash
# ${Math.random()}
exec rsync -a /root/pdns root@secondary:/root/
